# archive/deploy/bot.service
[Unit]
Description=Polymarket Arb Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket
ExecStart=/root/polymarket/venv/bin/uvicorn src.bot.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/root/polymarket/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-bot

[Install]
WantedBy=multi-user.target
